# jellyfin
#!/sbin/openrc-run

name=$RC_SVCNAME
description="Jellyfin Media Server"
supervisor="supervise-daemon"
command="/usr/bin/jellyfin"
command_args="$JELLYFIN_WEB_OPT $JELLYFIN_RESTART_OPT $JELLYFIN_FFMPEG_OPT $JELLYFIN_SERVICE_OPT $JELLYFIN_NOWEBAPP_OPT $JELLYFIN_ADDITIONAL_OPTS"
supervise_daemon_args=" -d /var/lib/jellyfin"
command_user="jellyfin:jellyfin"

depend() {
        after network-online 
}
